# hw08.systemd/restart-kafka.path
[Unit]
Description=Path to restart kafka

[Path]
PathExists=/home/vagrant/kafka/error.log
Unit=restart-kafka.service